# systemd/logai-avante.service
[Unit]
Description=Log AI usage when Avante prompts are created
After=network.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/logai.sh
StandardOutput=journal
StandardError=journal

# Optional: Add a delay to batch multiple quick writes
ExecStartPre=/bin/sleep 2

# Optional: Prevent the service from running too frequently
# Limit to once every 10 seconds
RestartSec=10

[Install]
WantedBy=default.target

